# deploy/systemd/testsdkflow-server.service
[Unit]
Description=testsdkflow Next.js mini-app
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/testsdkflow
EnvironmentFile=/opt/testsdkflow/.env
Environment=NODE_ENV=production
Environment=PORT=3010
ExecStart=/usr/bin/npx --no-install next start -p 3010
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
